# src/xarcade2jstick.service
[Unit]
Description=Xarcade2Jstick

[Service]
Type=oneshot
ExecStart=/usr/local/bin/xarcade2jstick
